# cron/crontab
* * * * * root echo " /code/test.py >> /var/log/cron.log
# Empty line
